# release/linux/fpm.before-remove.sh
systemctl stop safe-t-daemon.service
systemctl disable safe-t-daemon.service
