# release/linux/fpm.after-install.sh
systemctl enable safe-t-daemon.service
systemctl start safe-t-daemon.service
